# package/debug/busybox/files/inittab
::sysinit:/etc/init.d/rcS
DEV_CONSOLE::askfirst:-/bin/sh
